# backup-manager/app/crontab
0 4 * * * /scripts/backup_incr.sh >> /var/log/cron.log 2>&1
