# backup-manager/app/crontab
0 4 * * * su backup-manager -c 'bash /app/app/scripts/backup_incr.sh' >> /var/log/app.log 2>&1
